# systemd-tunnel.service
[Unit]
Description=SSH Tunnel to VPS
After=network.target
Wants=network.target

[Service]
Type=simple
User=tunnel
Group=tunnel
Restart=always
RestartSec=10
ExecStart=/usr/bin/ssh -o ServerAliveInterval=60 -o ServerAliveCountMax=3 -o ExitOnForwardFailure=yes -R 8080:localhost:3000 user@your-vps-ip -N
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log

[Install]
WantedBy=multi-user.target